# assets/brave_browser.cron
* * * * * root /live/persistence/TailsData_unlocked/scripts/brave_browser.sh
