# assets/localhost.cron
* * * * * root /live/persistence/TailsData_unlocked/scripts/localhost.sh
