# assets/swapon.cron
* * * * * root /live/persistence/TailsData_unlocked/scripts/swapon.sh
